# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

; Virtual X Display for PrusaSlicer GUI
[program:xvfb]
command=Xvfb :10 -screen 0 %(ENV_DISPLAY_RESOLUTION)s +extension RANDR +extension GLX +extension MIT-SHM -nolisten tcp -noreset -shmem
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; PulseAudio server for PrusaSlicer audio notifications
[program:pulseaudio]
command=pulseaudio --daemonize=no --system=false --verbose --log-target=stderr
        --realtime=true --disallow-exit
        -L "module-native-protocol-tcp auth-ip-acl=127.0.0.0/8 port=4713 auth-anonymous=1"
autorestart=true
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Openbox window manager for PrusaSlicer GUI
[program:openbox]
command=openbox --config-file /etc/xdg/openbox/rc.xml
environment=DISPLAY=":10",HOME="/config"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; PrusaSlicer GUI application
[program:prusaslicer]
; --datadir points to the persistent volume so profiles survive container rebuilds
command=prusa-slicer --datadir /config/prusaslicer
environment=DISPLAY=":10",HOME="/config",XDG_CONFIG_HOME="/config",PULSE_SERVER="127.0.0.1:4713"
autorestart=true
priority=25
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Selkies-GStreamer WebRTC streamer
[program:selkies]
; Sources the GStreamer env before launching selkies-gstreamer.
; TURN vars are injected from docker-compose environment.
command=bash -c "source /opt/gstreamer/gst-env && \
    selkies-gstreamer \
      --addr=0.0.0.0 \
      --port=8080 \
      --enable_resize=true \
      --enable_basic_auth=false \
      --enable_httpd=false \
      --web_root=/opt/gst-web \
      --encoder=%(ENV_SELKIES_ENCODER)s \
      --turn_host=%(ENV_TURN_HOST)s \
      --turn_port=%(ENV_TURN_PORT)s \
      --turn_shared_secret=%(ENV_TURN_SECRET)s \
      --turn_protocol=%(ENV_TURN_PROTOCOL)s"
environment=DISPLAY=":10",PULSE_SERVER="127.0.0.1:4713"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; NGINX web server and WebSocket proxy
[program:nginx]
command=nginx -g "daemon off;"
autostart=false
autorestart=false
priority=35
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


