# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=root

; 1. Virtual X display
[program:xvfb]
command=Xvfb :10 -screen 0 1920x1080x24 +extension RANDR +extension GLX +extension MIT-SHM -nolisten tcp -noreset -shmem
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 2. PulseAudio
[program:pulseaudio]
command=pulseaudio --daemonize=no --system=false --verbose --log-target=stderr
        --realtime=true --disallow-exit
        -L "module-native-protocol-tcp auth-ip-acl=127.0.0.0/8 port=4713 auth-anonymous=1"
autorestart=true
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 3. Openbox window manager 
[program:openbox]
command=openbox --config-file /etc/xdg/openbox/rc.xml
environment=DISPLAY=":10",HOME="/config"
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 4. PrusaSlicer
[program:prusaslicer]
command=prusa-slicer --datadir /config/prusaslicer
environment=DISPLAY=":10",HOME="/config",XDG_CONFIG_HOME="/config",PULSE_SERVER="127.0.0.1:4713",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
autorestart=true
priority=25
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; 5. Selkies-GStreamer WebRTC streamer 
; TURN vars are resolved from env by entrypoint.sh before supervisord starts.
[program:selkies]
command=bash -c "source /opt/gstreamer/gst-env && \
    selkies-gstreamer \
      --addr=0.0.0.0 \
      --port=8080 \
      --web_root=/opt/gst-web/src \
      --enable_basic_auth=false \
      --enable_https=false \
      --enable_resize=true \
      --encoder=SELKIES_ENCODER_VALUE \
      --turn_host=TURN_HOST_VALUE \
      --turn_port=TURN_PORT_VALUE \
      --turn_shared_secret=TURN_SECRET_VALUE \
      --turn_protocol=TURN_PROTOCOL_VALUE"
environment=DISPLAY=":10",PULSE_SERVER="127.0.0.1:4713",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
